# roles/home-assistant/files/etc/systemd/system/home-assistant.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/hass --config /var/lib/hass/ --log-rotate-days 1 --pid-file /var/lib/hass/hass.pid
